# second/htg2.service
[Unit]
Description=How-To Geek Sample Service

[Service]
ExecStart=/usr/local/bin/htg2.sh
#with ExecStart and nothing else, assumes "Type=simple"

